# images/base-session/config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:vnc]
command=/usr/local/bin/vnc_startup.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/vnc.err.log
stdout_logfile=/var/log/supervisor/vnc.out.log

[program:xfce]
command=/usr/bin/xfce4-session
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/xfce.err.log
stdout_logfile=/var/log/supervisor/xfce.out.log
environment=DISPLAY=":1"
